# scripts/pipeline/net.litbulb.pipeline.test.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>net.litbulb.pipeline.test</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>TEST_TOPIC="Claude Sonnet 5" TEST_CATEGORY="AI" TEST_SLUG="claude-sonnet-5-test" /Users/hus9/Claude/litbulb/abedis-lab/scripts/pipeline/run_pipeline.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Year</key>
    <integer>2026</integer>
    <key>Month</key>
    <integer>7</integer>
    <key>Day</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>19</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/hus9/Claude/litbulb/abedis-lab/scripts/pipeline/logs/launchd-test-out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/hus9/Claude/litbulb/abedis-lab/scripts/pipeline/logs/launchd-test-err.log</string>
</dict>
</plist>
